# pi_bullseye_backup/home/pi/.config/autostart/light-locker.desktop
[Desktop Entry]
Type=Application
Name=Screen Locker
Exec=light-locker
Hidden=true
X-GNOME-Autostart-enabled=false
